# wolapp.service
[Unit]
Description=wolapp - Home Control (gunicorn)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=foreigndev
WorkingDirectory=/home/foreigndev/wolapp
EnvironmentFile=/home/foreigndev/wolapp/.env
ExecStart=/bin/bash -lc 'BIND=$$(tailscale ip -4 2>/dev/null || echo 127.0.0.1); exec /home/foreigndev/wolapp/venv/bin/gunicorn --bind "$${BIND}:8080" --workers 1 --threads 2 --timeout 120 --access-logfile /dev/null --error-logfile - app:app'
Restart=always
RestartSec=3
RuntimeWatchdogSec=120
WatchdogSignal=SIGKILL
MemoryMax=200M

[Install]
WantedBy=multi-user.target
